# package/ahoy/ahoy/files/ahoy.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013  

START=40

start() {
  cd /ahoy
  start-stop-daemon -S -b -m -p /var/run/ahoy.pid -x /ahoy/ahoy.lua
}

reload() {
	return
}

stop() {
  start-stop-daemon -K -b -p /var/run/ahoy.pid 
}
